# mock-oracle-api.service
[Unit]
Description=Mock Oracle Database Middleware API
After=network.target

[Service]
Type=simple
User=stefanb
WorkingDirectory=/home/stefanb/oracle-operator-lab/mock-api
ExecStart=/home/stefanb/oracle-operator-lab/mock-api/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
